# ssdpBridge/arch/Denon.install
post_install() {
	systemctl enable ssdpBridge.service
}

pre_remove() {
	systemctl disable ssdpBridge.service
}
