# ssdpBridge/deb/prerm
#!/bin/sh
systemctl disable ssdpBridge.service
